# music-assistant-alexa-api/rootfs/etc/services.d/music-assistant-alexa-api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
set -e

# Declare variables
declare MA_HOSTNAME
declare API_PORT
declare API_USERNAME
declare API_PASSWORD
declare AWS_DEFAULT_REGION

MA_HOSTNAME=$(bashio::config 'ma_hostname' '')
API_PORT=$(bashio::addon.port '5000')
API_USERNAME=$(bashio::config 'api_username' '')
API_PASSWORD=$(bashio::config 'api_password' '')
AWS_DEFAULT_REGION=$(bashio::config 'aws_default_region' 'us-west-1')

if bashio::config.is_empty "api_password"; then
    API_PASSWORD=$(cat /dev/urandom | tr -dc 'A-Za-z0-9' | head -c 32) || true
    bashio::addon.option "api_password" "${API_PASSWORD}"
fi

bashio::log.info "Hostname: ${MA_HOSTNAME}"
bashio::log.info "Port: ${API_PORT}"
bashio::log.info "Username: ${API_USERNAME}"
bashio::log.info "Password: ${API_PASSWORD}"
bashio::log.info "AWS Region: ${AWS_DEFAULT_REGION}"

export MA_HOSTNAME
export API_PORT
export API_USERNAME
export API_PASSWORD
export AWS_DEFAULT_REGION

## Run your program
exec /app/venv/bin/gunicorn -w 1 app:app --chdir /app/src/ --reuse-port --bind 0.0.0.0:5000
